# crontab
*/15 * * * * /home/app-user/.local/share/virtualenvs/kyutechAppServer2018-yIFX9ipV/bin/python /home/app-user/kyutechAppServer2018/manage.py scrape_news >> ~/log/cron.log 2>> ~/log/scrape.log
