# crontab
*/15 * * * * /usr/bin/.pyenv/shims/python /home/ec2-user/kyutechAppServer2018/manage.py scrape_news >> ~/log/cron.log 2>> ~/log/test.log
